# dist/rinertia.service
[Unit]
Description=rinertia — inertial scrolling daemon for Linux touchpads
Documentation=https://github.com/JimMoen/rinertia
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/rinertia
Restart=on-failure
RestartSec=3

# To use a config file, override ExecStart:
#   systemctl --user edit rinertia.service
#   [Service]
#   ExecStart=
#   ExecStart=/usr/bin/rinertia --config %h/.config/rinertia/config.toml

[Install]
WantedBy=graphical-session.target
